# .config/systemd/user/vdirsyncer-sync.service
[Unit]
Description=Run vdirsyncer sync

[Service]
Type=oneshot
ExecStart=/usr/bin/vdirsyncer sync
